# data/com.mob.lotab.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mob.lotab</string>
    <key>ProgramArguments</key>
    <array>
        <string>@BINDIR@/daemon</string>
        <string>--app-path</string>
        <string>@BINDIR@/@APP_NAME@</string>
    </array>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/lotab.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/lotab.stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
